# linux_security/1x04_security_automation/sentinel.timer
[Unit]
Description=Run sentinel every 5 minute

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
